# deploy/trading-bot.service
[Unit]
Description=Autonomous Trading Bot
After=network-online.target
Wants=network-online.target
# Restart if the server reboots or the process crashes

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading-bot
EnvironmentFile=/opt/trading-bot/.env
ExecStart=/opt/trading-bot/venv/bin/python scheduler/main_loop.py

# Auto-restart on crash
Restart=always
RestartSec=30

# Logs go to journald (view with: journalctl -u trading-bot -f)
StandardOutput=journal
StandardError=journal

# Give the process time to shut down gracefully (sends SIGTERM first)
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
